# krep_tray.service
[Unit]
Description=Krep Tray App

[Service]
ExecStart=%h/.cargo/bin/krep-tray
Restart=on-failure

[Install]
WantedBy=default.target
